# logwatch.cron
MAILTO=root

04 1 * * * root /usr/sbin/logwatch-cron

# vim:ft=crontab
